# systemd/niriidle.service.in
[Unit]
Description=niri idle daemon
PartOf=graphical-session.target
After=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/niriidle
Restart=on-failure

[Install]
WantedBy=graphical-session.target
